# src/files/etc/hotplug.d/iface/99-custom
#!/bin/sh
[ "$ACTION" = "ifup" -a "$INTERFACE" = "wg0" ] && {
    logger "Fix IP Route"
    wg_scripts.sh fix
    # Add your custom command(s) here
}
